# src/include/apkvol/etc/apk/protected_paths.d/lbu.list
+etc/init.d/macaddress
+root/.ssh/authorized_keys
